# app-emulation/hyperv-integration-services/files/hyperv-integration-services.service.in
[Unit]
Description=Hyper-V @daemon@ Daemon
ConditionVirtualization=microsoft

[Service]
Type=forking
ExecStart=/usr/sbin/hv_@daemon@_daemon

[Install]
WantedBy=multi-user.target
