# app-emulation/hyperv-integration-services/files/hyperv-integration-services.service.in
[Unit]
Description=Hyper-V @daemon@ Daemon
ConditionVirtualization=microsoft

[Service]
ExecStart=/usr/sbin/hv_@daemon@_daemon --no-daemon

[Install]
WantedBy=multi-user.target
