# app-emulation/hyperv-integration-services/files/hyperv-integration-services.initd.in
#!/sbin/openrc-run
command="/usr/sbin/hv_@daemon@_daemon"
description="Hyper-V @daemon@ Daemon"

depend() {
	need localmount
}
